# ansible-collections/grnet/jitsi/roles/jitsi_excalidraw_backend/files/jitsi-excalidraw-backend.service
[Unit]
Description=Jitsi Excalidraw Backend

[Service]
User=excalidraw
ExecStart=npm start
WorkingDirectory=/opt/jitsi-excalidraw-backend

[Install]
WantedBy=multi-user.target
